# PATCH/files/etc/hotplug.d/net/01-maximize_nic_rx_tx_buffers
#!/bin/sh

[ "$ACTION" = "add" ] || exit 0
[ -z "$DEVICENAME" ] && exit 0

command -v ethtool >/dev/null 2>&1 || exit 0

NIC="$DEVICENAME"

[ ! -d "/sys/class/net/$NIC/device" ] && exit 0
[ "$(cat "/sys/class/net/$NIC/type" 2>/dev/null)" != "1" ] && exit 0

ETHTOOL_G=$(ethtool -g "$NIC" 2>/dev/null)
[ -z "$ETHTOOL_G" ] && exit 0

RX_MAX=$(echo "$ETHTOOL_G" | awk '/RX:/ {print $2; exit}')
TX_MAX=$(echo "$ETHTOOL_G" | awk '/TX:/ {print $2; exit}')

[ -n "$RX_MAX" ] && [ "$RX_MAX" != "0" ] && ethtool -G "$NIC" rx "$RX_MAX" 2>/dev/null
[ -n "$TX_MAX" ] && [ "$TX_MAX" != "0" ] && ethtool -G "$NIC" tx "$TX_MAX" 2>/dev/null

if ethtool -c "$NIC" 2>/dev/null | grep -q "Adaptive RX:"; then
    ethtool -C "$NIC" adaptive-rx on adaptive-tx on >/dev/null 2>&1
fi

logger -t net-buffer "Maximized ring buffers for $NIC"
