# PATCH/files/etc/hotplug.d/iface/99-zzz-odhcpd
#!/bin/sh

[ "$ACTION" = "ifup" ] || [ "$ACTION" = "ifupdate" ] || exit 0
[ "$ACTION" = "ifupdate" ] && [ -z "$IFUPDATE_ADDRESSES" ] && [ -z "$IFUPDATE_DATA" ] && exit 0

. /lib/functions/network.sh
network_flush_cache
network_find_wan6 NET_IF6
[ "$INTERFACE" = "$NET_IF6" ] || exit 0

DEBOUNCE_FILE="/tmp/run/odhcpd_hotplug_gen"
mkdir -p /tmp/run

my_gen="$$.$(date +%s)"
echo "$my_gen" > "$DEBOUNCE_FILE"

(
    sleep 5
    [ "$(cat "$DEBOUNCE_FILE" 2>/dev/null)" = "$my_gen" ] || exit 0

    PING6=$(command -v ping6 || echo "ping -6")

    uci -q set dhcp.lan.ra_mininterval=1
    uci -q set dhcp.lan.ra_maxinterval=2
    uci -q set dhcp.lan.ra_lifetime=0
    uci -q set dhcp.lan.ra_preferred_lifetime=0
    uci -q set dhcp.lan.ra_valid_lifetime=0
    /etc/init.d/odhcpd restart >/dev/null 2>&1

    network_get_device LAN_DEVICE lan
    for dev in ${LAN_DEVICE:-br-lan}; do
        $PING6 -c 2 -I "$dev" ff02::1 >/dev/null 2>&1 &
    done

    sleep 4
    [ "$(cat "$DEBOUNCE_FILE" 2>/dev/null)" = "$my_gen" ] || exit 0

    uci -q revert dhcp.lan.ra_lifetime
    uci -q revert dhcp.lan.ra_preferred_lifetime
    uci -q revert dhcp.lan.ra_valid_lifetime
    uci -q set dhcp.lan.ra_mininterval=1
    uci -q set dhcp.lan.ra_maxinterval=3
    /etc/init.d/odhcpd restart >/dev/null 2>&1

    for dev in ${LAN_DEVICE:-br-lan}; do
        $PING6 -c 2 -I "$dev" ff02::1 >/dev/null 2>&1 &
    done

    sleep 10
    [ "$(cat "$DEBOUNCE_FILE" 2>/dev/null)" = "$my_gen" ] || exit 0

    uci -q revert dhcp.lan.ra_mininterval
    uci -q revert dhcp.lan.ra_maxinterval
    /etc/init.d/odhcpd restart >/dev/null 2>&1

    [ -x /etc/init.d/miniupnpd ] && /etc/init.d/miniupnpd restart >/dev/null 2>&1

    logger -t odhcpd "IPv6 elegant refresh finished for $INTERFACE"
) &
